# relay/wireless-mic-relay.service
[Unit]
Description=Wireless Mic Input Relay
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/wireless-mic-relay
ExecStart=/usr/bin/node relay.js
Restart=always
RestartSec=5
Environment=PORT=4001
# Uncomment and set these as needed:
# Environment=RELAY_SECRET=your-strong-secret-here
# Environment=CERT_DIR=/etc/letsencrypt/live/yourdomain.com
# Environment=PUBLIC=/opt/wireless-mic-relay/public

[Install]
WantedBy=multi-user.target
